# server/wallet-app.service
[Unit]
Description=Home Accounting Legendary Django App
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=notify
User=walletapp
Group=walletapp
WorkingDirectory=/var/www/wallet.my-bucket.ru/home-accounting-legendary
Environment=PATH=/var/www/wallet.my-bucket.ru/venv/bin
Environment=DJANGO_SETTINGS_MODULE=core.settings.production
ExecStart=/var/www/wallet.my-bucket.ru/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 --timeout 120 --keep-alive 2 --max-requests 1000 --max-requests-jitter 100 core.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wallet-app

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/wallet.my-bucket.ru
ReadWritePaths=/var/log/django

[Install]
WantedBy=multi-user.target
